# usyslogd.service
[Service]
Restart=on-success
ExecStart=/home/snide/src/usyslogd/usyslogd
