# usyslogd.service
[Service]
Restart=on-success
ExecStart=/usr/local/sbin/usyslogd
SyslogIdentifier=usyslogd

[Install]
WantedBy=multi-user.target
